# .config/systemd/user/low_bat.service
[Unit]
Description=check for low battery

[Service]
ExecStart=bash /usr/local/bin/low_bat

[Install]
WantedBy=multi-user.target
